# deploy/systemd/pathonet-hub.service
[Unit]
Description=PathoNet Mock Hub
After=network-online.target
Wants=network-online.target
PartOf=pathonet.target

[Service]
Type=notify
NotifyAccess=main
Environment=PATHONET_ROOT=/opt/pathonet/current
EnvironmentFile=-/etc/pathonet/pathonet.env
WorkingDirectory=${PATHONET_ROOT}/hub
ExecStart=/usr/bin/dotnet ${PATHONET_ROOT}/hub/PathoNet.Hub.dll
Restart=always
RestartSec=5
WatchdogSec=30
KillSignal=SIGTERM
TimeoutStopSec=20
SyslogIdentifier=pathonet-hub

[Install]
WantedBy=pathonet.target
