# deploy/systemd/pathonet-apisender.service
[Unit]
Description=PathoNet Mock ApiSender
After=pathonet-api.service pathonet-hub.service network-online.target
Wants=pathonet-api.service pathonet-hub.service network-online.target
PartOf=pathonet.target

[Service]
Type=notify
NotifyAccess=main
Environment=PATHONET_ROOT=/opt/pathonet/current
EnvironmentFile=-/etc/pathonet/pathonet.env
WorkingDirectory=${PATHONET_ROOT}/apisender
ExecStart=/usr/bin/dotnet ${PATHONET_ROOT}/apisender/PathoNet.ApiSender.dll
Restart=always
RestartSec=5
WatchdogSec=30
KillSignal=SIGTERM
TimeoutStopSec=20
SyslogIdentifier=pathonet-apisender

[Install]
WantedBy=pathonet.target
